# p44bandit-config/files/etc/service/p44bandit/run
#!/bin/sh

# default features
EXECUTABLE=p44banditd
LOGLEVEL=5
EXTRAOPTS=
MESSAGE=


# run the daemon
exec ${EXECUTABLE} \
  --jsonapiport 8090 \
  -l ${LOGLEVEL} \
  ${EXTRAOPTS} \
  </dev/null
